# images/marathon/slack-relay/resources/supervisor/pod.conf
[program:relay]
command=python /opt/relay/pod/pod.py