# images/marathon/haproxy/resources/supervisor/haproxy.conf
[program:haproxy]
command=python /opt/haproxy/pod/pod.py
