# images/marathon/slave/resources/supervisor/slave.conf
[program:slave]
command=python /opt/slave/pod/pod.py
