# images/marathon/slave/resources/supervisor/socat.conf
[program:socat]
command=socat -d -d TCP-LISTEN:9001,fork UNIX:/host/var/run/docker.sock
