# images/marathon/hook/resources/supervisor/hook.conf
[program:hook]
command=python /opt/hook/pod/pod.py
