# images/marathon/redis/resources/supervisor/redis.conf
[program:redis]
command=python /opt/redis-2.8.17/pod/pod.py
