# images/marathon/servo/resources/supervisor/servo.conf
[program:servo]
command=python /opt/servo/pod/pod.py
